# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log

[program:node]
command=node /usr/src/app/index.js

[program:java]
command=java -Djava.security.egd=file:/dev/./urandom -jar /app.jar